# supervisord/conf.d/gunicorn.conf

[program:gunicorn]
command=/Users/irio/Desktop/studying/drf/drf_booking/venv/bin/gunicorn -c /Users/irio/Desktop/studying/drf/drf_booking/gunicorn/dev.py
directory=/Users/irio/Desktop/studying/drf/drf_booking/config
user=nobody
autostart=true
autorestart=true
logfile = Users/irio/Desktop/studying/drf/drf_booking/logs/supervisord.log
