# supervisord/conf.d/gunicorn.conf

[program:gunicorn]
command=/Users/irio/Desktop/studying/drf/drf_booking/venv/bin/gunicorn -c /Users/irio/Desktop/studying/drf/drf_booking/gunicorn/dev.conf.py
directory=/Users/irio/Desktop/studying/drf/drf_booking/config
user=root
autostart=true
autorestart=true
startsecs=0
childlogdir=/Users/irio/Desktop/studying/drf/drf_booking/logs